# dotfiles/clipmenu/.config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment="CM_IGNORE_WINDOW=KeePassXC"
Environment="CM_MAX_CLIPS=5"
